# deploy/paideia-backup.timer
[Unit]
Description=Run Paideia backup daily at 04:00 UTC

[Timer]
OnCalendar=*-*-* 04:00:00 UTC
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
